# etc-init-conf/dubbo.conf
start on dub
stop on runlevel [016]

respawn
respawn limit 1 5

script
  #. /etc/profile.d/caiex.sh
  exec /opt/rclocal/dubbo_upstart.sh
end script

