# etc-init-conf/matlab.conf
start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit 1 5

script
  . /etc/profile.d/caiex.sh

  exec \
  /opt/rclocal/matlab_service_upstart.sh \
end script


