# etc-init-conf/sbcconsumer.conf
start on sbccon
#start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit 1 5

script
  . /etc/profile.d/caiex.sh
  /opt/rclocal/sbcconsumer_upstart.sh

end script
