# systemd/1.5/consumer.service
[Unit]
Description=OLTP-consumer
After=network.target
After=syslog.target
#Requires=sbcconsumer.service

[Install]
WantedBy=multi-user.target

[Service]
User=root
Group=caiex


# Start main service
ExecStart=/opt/rclocal/consumer_upstart.sh


# Give up if ping don't get an answer
TimeoutSec=600

Restart=always
PrivateTmp=false
